# coolify/rootfs/etc/services.d/coolify/finish
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Cleanup after Coolify stops
# ==============================================================================

bashio::log.info "Coolify stopped"
